# files/etc/udev/rules.d/99-usbcopy.rules
KERNEL=="sd?", SUBSYSTEMS=="usb", ACTION=="add", RUN+="/usr/local/sbin/copy_img $kernel"
