# files/etc/udev/rules.d/99-usbcopy.rules
KERNEL=="sd?", SUBSYSTEMS=="usb", ACTION=="add", RUN+="/bin/systemctl --no-block start copy_img2@$kernel.service"
